# debian/debusine-server.debusine-server-delete-expired.timer
[Unit]
Description=Daily cleanup of Debusine expired or unused objects

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
